# search-deb/src/deb/control/prerm
#!/bin/sh
systemctl stop    search || true
systemctl disable search || true
